# root_forked-daapd/etc/services.d/dbus/run
#!/usr/bin/with-contenv bash

echo "Start service: dbus-daemon"
exec \
	dbus-daemon --system --nofork
