# root_forked-daapd/etc/services.d/eudev/run
#!/usr/bin/with-contenv bash

echo "Start service: eudev"
exec \
	/sbin/udevd
